# deploy/qbot-komoot-watch.service
[Unit]
Description=QBot Komoot watcher - wykrywa nowe/zmienione trasy Komoot i pyta na Telegramie (bez auto-pushu)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=qbot
Group=qbot
WorkingDirectory=/opt/qbot/app
ExecStart=/opt/qbot/app/.venv/bin/python /opt/qbot/app/komoot_watch.py
TimeoutStartSec=300
